# debian/tty-egpf-monitord.prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
  if command -v systemctl >/dev/null 2>&1; then
    systemctl stop tty-egpf-monitord.service || true
    systemctl disable tty-egpf-monitord.service || true
  fi
fi

#DEBHELPER#


