# packaging/tty-egpf-monitord.service
[Unit]
Description=TTY eBPF monitor daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tty-egpf-monitord --socket /run/tty-egpf-monitord.sock --log-dir /var/log/tty-egpf-monitor
Restart=on-failure
User=root
Group=root
RuntimeDirectory=tty-egpf-monitord
NoNewPrivileges=false
ProtectSystem=no
ProtectHome=no
ReadWritePaths=/var/log/tty-egpf-monitor /run /sys/kernel/debug
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_NET_ADMIN CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_NET_ADMIN CAP_SYS_ADMIN
LimitMEMLOCK=infinity
RestartSec=2s
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target


